# serverService/gwentapi.socket
[Unit]
Description=GwentAPI Socket
PartOf=gwentapi.service

[Socket]
ListenStream=127.0.0.1:8080

[Install]
WantedBy=sockets.target